# raspireader/raspireader-update.service
[Unit]
Description=RaspiReader OTA Update

[Service]
Type=oneshot
ExecStart=/opt/raspireader/repo/raspireader/update.sh
